# services/axi-lab.service
[Unit]
Description=AxiDraw HTTP Server
Documentation=https://github.com/joeysapp/axi-lab
After=network.target

[Service]
Type=simple
User=zooey
Group=staff
WorkingDirectory=/Users/zooey/etc/git/axi-lab
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10

# Environment variables
Environment=AXIDRAW_PORT=9700
Environment=AXIDRAW_AUTO_CONNECT=true
Environment=AXIDRAW_NARROW_BAND=true
Environment=AXIDRAW_MODEL=V2_V3
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=axi-lab

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
