# services/axidraw-tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.axidraw.tmux</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            # Kill any existing axi session
            /usr/local/bin/tmux kill-session -t axi 2>/dev/null || true

            # Start new tmux session with the server
            /usr/local/bin/tmux new-session -d -s axi -c "$HOME/etc/git/axi-server"

            # Set environment and start server
            /usr/local/bin/tmux send-keys -t axi 'export AXIDRAW_PORT=9700' Enter
            /usr/local/bin/tmux send-keys -t axi 'export AXIDRAW_AUTO_CONNECT=true' Enter
            /usr/local/bin/tmux send-keys -t axi 'node src/index.js' Enter
        </string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/zooey</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/axidraw-tmux.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/axidraw-tmux.error.log</string>
</dict>
</plist>
