# services/com.axidraw.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.axidraw.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd $HOME/etc/git/axi-server &amp;&amp; exec /usr/local/bin/node src/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AXIDRAW_PORT</key>
        <string>9700</string>
        <key>AXIDRAW_AUTO_CONNECT</key>
        <string>true</string>
        <key>AXIDRAW_NARROW_BAND</key>
        <string>true</string>
        <key>AXIDRAW_MODEL</key>
        <string>SE_A3</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/axidraw-server.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/axidraw-server.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/zooey/etc/git/axi-server</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
